# roles/internal/planningalerts/files/planningalerts-worker.conf
start on starting planningalerts
stop on stopping planningalerts
